# deploy/supervisor/noiachat-worker.conf
[program:noiachat-worker]
process_name=%(program_name)s_%(process_num)02d
command=php8.4 /var/www/noiachat/artisan queue:work --sleep=3 --tries=3 --timeout=120 --max-time=3600
directory=/var/www/noiachat
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/noiachat/storage/logs/worker.log
stopwaitsecs=3600
